# deploy/event-tool-backend.service
[Unit]
Description=Event Tool Backend (FastAPI on :8002)
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/event-tool/backend
Environment="PATH=/var/www/event-tool/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/var/www/event-tool/backend/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8002
Restart=always
RestartSec=5
StandardOutput=append:/var/log/event-tool-backend.log
StandardError=append:/var/log/event-tool-backend.log

[Install]
WantedBy=multi-user.target
